# debian/lvmsnapper.cron.hourly
#!/bin/sh

ENABLE="1"
SNAPPER=$(which lvmsnapper.py)


if [ "$ENABLE" = '1' ]
then
	$SNAPPER
fi
