# debian/lvmsnapper.cron.hourly
#!/bin/sh

ENABLE=
SNAPPER=$(which lvmsnapper)


if [ "$ENABLE" = '1' ]
then
	$SNAPPER
fi
